# services/backup.timer
[Unit]
Description=Run backup every night

[Timer]
OnCalendar=*-*-* 00:00:00
Persistent=true

[Install]
WantedBy=timers.target
